# systemd/account-watch.timer
[Unit]
Description=Run account-watch every 5 minutes (autonomous account failover)

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
